# Docker_JuypterLab/config/supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root
autorestart=true